# assets/systemd/coin-price.timer
[Unit]
Description=Update crpto currency prices every hour

[Timer]
OnCalendar=hourly
Unit=coin-price.service

[Install]
WantedBy=timers.target
